# update.crontab
0 * * * * bash -c 'cd ~/frontend-mock && ./update.sh' >~/frontend-mock/update.txt 2>&1
