# update.crontab
0 * * * * bash -c 'cd ~/api && ./update.sh' >~/api/update.txt 2>&1
